# infra/systemd/tpb-data-enrichment.service
[Unit]
Description=TPB Data-enrichment Service
After=network.target postgresql.service redis.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/project/AgentdesktopTPB
Environment="NODE_ENV=production"
Environment="PORT=3016"
ExecStart=/usr/bin/npx nx run data-enrichment-service:serve
Restart=always
RestartSec=10
StandardOutput=append:/opt/project/AgentdesktopTPB/logs/data-enrichment-service.log
StandardError=append:/opt/project/AgentdesktopTPB/logs/data-enrichment-service.log

[Install]
WantedBy=multi-user.target
