# infra/systemd/tpb-notification.service
[Unit]
Description=TPB Notification Service
After=network.target postgresql.service redis.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/project/AgentdesktopTPB
Environment="NODE_ENV=production"
Environment="PORT=3006"
ExecStart=/usr/bin/npx nx run notification-service:serve
Restart=always
RestartSec=10
StandardOutput=append:/opt/project/AgentdesktopTPB/logs/notification-service.log
StandardError=append:/opt/project/AgentdesktopTPB/logs/notification-service.log

[Install]
WantedBy=multi-user.target
