# infra/systemd/tpb-dashboard.service
[Unit]
Description=TPB Dashboard Service
After=network.target postgresql.service redis.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/project/AgentdesktopTPB
Environment="NODE_ENV=production"
Environment="PORT=3017"
ExecStart=/usr/bin/npx nx run dashboard-service:serve
Restart=always
RestartSec=10
StandardOutput=append:/opt/project/AgentdesktopTPB/logs/dashboard-service.log
StandardError=append:/opt/project/AgentdesktopTPB/logs/dashboard-service.log

[Install]
WantedBy=multi-user.target
