# infra/systemd/tpb-object-schema.service
[Unit]
Description=TPB Object-schema Service
After=network.target postgresql.service redis.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/project/AgentdesktopTPB
Environment="NODE_ENV=production"
Environment="PORT=3013"
ExecStart=/usr/bin/npx nx run object-schema-service:serve
Restart=always
RestartSec=10
StandardOutput=append:/opt/project/AgentdesktopTPB/logs/object-schema-service.log
StandardError=append:/opt/project/AgentdesktopTPB/logs/object-schema-service.log

[Install]
WantedBy=multi-user.target
